# provisioning/ansible/roles/external_mail_handler/templates/externalmailhandler.service.j2
[Unit]
Description=SOCBED External Mail Handler
After=network.target

[Service]
Type=simple
ExecStart={{ exec_cmd }}
Restart=always

[Install]
WantedBy=multi-user.target
